# crontab
*/5 * * * * regfish_dyndns2_start.sh >> /dev/stdout 2>&1
# crontab requires an empty line at the end of the file
